# src/systemd/opencode-metrics.service
[Unit]
Description=OpenCode Metrics Exporter (Prometheus)
After=network.target opencode-infra.service
Wants=network.target opencode-infra.service

[Service]
Type=simple
Environment=PATH=%h/.local/bin:%h/.n/bin:%h/.bun/bin:/usr/local/bin:/usr/bin:/bin
Environment=METRICS_EXPORTER_PORT=9464
Environment=DEPLOYMENT_PROFILE=personal
ExecStart=python3 %h/Projects/opencode_initializer/scripts/oc-metrics.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
